# files/usr/lib/systemd/system/bbsd.service
[Unit]
Description=bbsd

[Service]
ExecStart=/usr/local/bbsd/pmain
TasksMax=infinity

[Install]
WantedBy=multi-user.target
Alias=bbsd.service
